# all-in-one/supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── 1. onec-toolkit (1C live data MCP server) ─────────────────────────────────
[program:onec-toolkit]
command=python -m onec_mcp_toolkit_proxy
directory=/app
environment=PORT="6003",TIMEOUT="%(ENV_ONEC_TIMEOUT)s",RESPONSE_FORMAT="json",ALLOW_DANGEROUS_WITH_APPROVAL="true"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/onec-toolkit.log
stderr_logfile=/var/log/supervisor/onec-toolkit.log
stdout_logfile_maxbytes=5MB

# ── 2. BSL Language Server session manager ────────────────────────────────────
[program:bsl-session-manager]
command=/usr/bin/lsp-session-manager
    --port=9999
    --workspace=%(ENV_BSL_WORKSPACE)s
    --command=java
    --
    -Xmx%(ENV_BSL_JAVA_XMX)s
    -Xms%(ENV_BSL_JAVA_XMS)s
    -XX:+UseG1GC
    -XX:MaxGCPauseMillis=200
    -jar /opt/bsl-ls/bsl-language-server.jar
    lsp
    -c /etc/mcp-lsp-bridge/bsl-ls.json
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/bsl-session-manager.log
stderr_logfile=/var/log/supervisor/bsl-session-manager.log
stdout_logfile_maxbytes=5MB

# ── 3. platform-context (1C platform API docs) ────────────────────────────────
[program:platform-context]
command=java -jar /opt/platform/mcp-bsl-context.jar
    --mode sse
    --platform-path %(ENV_PLATFORM_PATH)s
    --port 8081
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/platform-context.log
stderr_logfile=/var/log/supervisor/platform-context.log
stdout_logfile_maxbytes=5MB

# ── 4. MCP Gateway (aggregator) ───────────────────────────────────────────────
[program:gateway]
command=python -m gateway
directory=/opt/gateway
environment=
    ONEC_TOOLKIT_URL="http://localhost:6003/mcp",
    PLATFORM_CONTEXT_URL="http://localhost:8081/sse",
    BSL_LSP_COMMAND="/usr/bin/mcp-lsp-bridge",
    LOG_LEVEL="%(ENV_LOG_LEVEL)s",
    ENABLED_BACKENDS="%(ENV_ENABLED_BACKENDS)s"
autostart=true
autorestart=true
startretries=3
startsecs=15
stdout_logfile=/var/log/supervisor/gateway.log
stderr_logfile=/var/log/supervisor/gateway.log
stdout_logfile_maxbytes=5MB
