# docker/lsp-bridge/s6-rc.d/bsl-ls/finish
#!/bin/sh
echo "BSL Language Server stopped with exit code $1"
